# init.d/botd.service
[Unit]
Description=botd

[Service]
ExecStart=/usr/local/bin/botd 

[Install]
WantedBy=multi-user.target
